# crontab
#crontab
* * * * * /usr/local/bin/python /app/run_app.py >> /app/cron.log 2>&1    
